# root/etc/services.d/apache2/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/httpd -DNO_DETACH
